# package/base-files/files/etc/init.d/eneth0
#!/bin/sh /etc/rc.common
START=98
start(){
	widora_mode router
}
