# package/base-files/files/etc/init.d/gprs
#!/bin/sh /etc/rc.common
START=99
start(){
	/etc/ppp/ppp-on
}
stop(){
        /etc/ppp/ppp-off
}
restart(){
       /etc/ppp/ppp-on
}
