# package/base-files/files/etc/init.d/lora
#!/bin/sh /etc/rc.common
START=99
start(){
        /bin/node-manage
}
stop(){
        /bin/node-manage -s
}
restart(){
       /bin/node-manage
}
