# package/base-files/files/etc/init.d/lora
#!/bin/sh /etc/rc.common
START=99
NAME=node-manage
start() {
	/bin/node-manage&
}

stop() {
	killall node-manage
}
restart() {
	killall node-manage
	/bin/node-manage&
}
